# Empacotamento/mueller/mueller-2002.02.27/debian/mueller7-dict/DEBIAN/postrm
#!/bin/sh



set -e

if [ "$1" = "remove" ] || [ "$1" = "purge" ]; then
  dictdconfig -w 2>/dev/null || :
  invoke-rc.d dictd force-reload 2>/dev/null || :
fi
